# deploy/systemd/zhoujinxin-portfolio-cert-renew.service
[Unit]
Description=Renew the Zhou Jinxin portfolio IP certificate
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/sh /srv/zhoujinxin-portfolio/scripts/renew-ip-certificate.sh
